# board/tind/compact-rtu/overlay/etc/systemd/system/cron.service
[Unit]
Description=Command Scheduler
RequiresMountsFor=/var/spool/cron/

[Service]
ExecStart=/usr/sbin/crond
KillMode=process
Type=forking

[Install]
WantedBy=multi-user.target

