# board/tind/compact-rtu/overlay/etc/systemd/system/cron.service
[Unit]
Description=Command Scheduler

[Service]
ExecStart=/usr/sbin/crond
KillMode=process

[Install]
WantedBy=multi-user.target

