# board/tind/compact-rtu/overlay/etc/systemd/system/print-ppp0.service
[Unit]
Description=Print modem status to console
DefaultDependencies=yes

[Service]
Type=oneshot
ExecStart=/usr/sbin/modem-status

