# board/suzipanel/overlay/etc/profile.d/opendaf.sh
export OPENDAF_CFGROOT=/etc/opendaf
export OPENDAF_PREFIX=/usr

