# board/tind/compact-rtu/mod-wirt/udev/90-rc232.rules
SUBSYSTEM=="gpio", KERNEL=="gpio71", RUN+="/bin/ln -sf /sys/%p /dev/rc232_config"
SUBSYSTEM=="tty", KERNEL=="ttyAPP1", SYMLINK+="ttyRC232"

