# board/ha2/overlay/etc/profile.d/export_mac.sh
export MAC=`cat /sys/class/net/eth0/address | sed s/://g | tr 'a-z' 'A-Z'`
